# cookbooks/docker/files/etc/systemd/system/docker-prune.service
[Unit]
Description=Docker system prune
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune --force
